# modules/post-install/atlas-snapshot/atlas-snapshot.timer
[Unit]
Description=Daily timer for Atlas snapshots

[Timer]
OnCalendar=00:00:00
Persistent=true

[Install]
WantedBy=timers.target
