# modules/post-install/atlas-snapshot/atlas-snapshot.service
[Unit]
Description=Atlas storage snapshot backup

[Service]
Type=oneshot
User=__SYSTEM_USER__
ExecStart=/usr/local/bin/atlas-snapshot
StandardOutput=journal
StandardError=journal
