# modules/atlas/atlas-snapshot.service
[Unit]
Description=Atlas storage snapshot backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/atlas-snapshot
StandardOutput=journal
StandardError=journal
